# gunicorn.conf
[program:gunicorn]
command=gunicorn --workers=3 wsgi:app -b 0.0.0.0:5000
directory= /deploy
